# static_confs/etc/systemd/system/getty@tty1.service.d/override.conf
[Unit]
Description=Autostarts desktop
After=httpd.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin http --noclear %I $TERM
Type=simple
